# deploy/files/etc/systemd/system/serial-getty@ttyGS0.service.d/override.conf
[Service]
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=no

